# nginx-php/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm

[program:nginx]
command=/usr/sbin/nginx

[program:memcached]
command=/usr/bin/memcached
user=memcached
